# m/coder/coder.service
[Unit]
Description=Coder server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/home/ubuntu/m/coder/coder.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

